# portal/deploy/files/supervisor_gunicorn.conf
[program:gunicorn]
command={{ python }} {{project_path }}manage.py run_gunicorn -c {{ project_path }}gunicorn.conf.py
user={{ user }}
redirect_stderr=True
directory={{ project_path }}

